# root/etc/uci-defaults/97_luci
uci -q batch << EOI
set luci.main.lang='en'
commit luci
EOI
