# root/etc/uci-defaults/99_firewall
uci -q batch << EOI
add firewall rule
set firewall.@rule[9]=rule
set firewall.@rule[9].name='Allow-LuCi'
set firewall.@rule[9].src='*'
set firewall.@rule[9].dest_port='80'
set firewall.@rule[9].proto='tcp'
set firewall.@rule[9].target='ACCEPT'
add firewall rule
set firewall.@rule[10]=rule
set firewall.@rule[10].name='Allow-SSH'
set firewall.@rule[10].src='*'
set firewall.@rule[10].dest_port='22'
set firewall.@rule[10].target='ACCEPT'
commit firewall
EOI
